# roles/demo-uploader/files/upload-demo.service
[Unit]
Description=Upload CSGO demoes

[Service]
Type=oneshot
ExecStart=/home/steam/upload-demoes.sh
User=steam
Group=steam
