# roles/demo-uploader/files/upload-demo.timer
[Unit]
Description=Upload CS:GO demoes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=upload-demo.service

[Install]
WantedBy=timers.target
